# rootfs/etc/init.d/timestamp_updater
#!/sbin/openrc-run

PIDFILE="/var/run/timestamp_updater.pid"

depend() {
    need rclone
}

start() {
    ebegin "Starting timestamp updater"

    start-stop-daemon --start --background --pidfile "$PIDFILE" --make-pidfile --exec /usr/local/bin/timestamp_updater.sh

    eend $?
}

stop() {
    ebegin "Stopping timestamp updater"

    start-stop-daemon --stop --pidfile "$PIDFILE"

    eend $?
}

status() {
    pid=$(cat $PIDFILE)

    if [ -n "$pid" ]; then
        eend 0
    else
        eend 1
    fi
}
