# rootfs/etc/s6-overlay/s6-rc.d/rclone-serve/run
#!/usr/bin/env sh

setsid /usr/bin/rclone rcd \
    --config=/etc/rclone/rclone.conf \
    --rc-addr=:5572 \
    --rc-allow-origin="*" \
    --rc-files=/var/rclone/webgui \
    --rc-no-auth \
    --rc-web-gui-no-open-browser
