# rootfs/etc/s6-overlay/s6-rc.d/nfs-configure/up
#!/usr/bin/env sh

exportfs -r

rpc.nfsd -N 3
